# bot/spikelines-bot.service
# systemd unit for the Spikelines Telegram bot on the box.
# Install: sudo cp spikelines-bot.service /etc/systemd/system/ && sudo systemctl daemon-reload
#          sudo systemctl enable --now spikelines-bot
#
# Hard resource caps so the bot can NEVER starve the lagisalpha crons sharing this
# box, and auto-restart so a crash mid-match self-heals (all call state is in the
# DB, so a restart resumes cleanly).
[Unit]
Description=Spikelines Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs from inside the repo checkout so the bot can import ../lib (match-feed, markets).
WorkingDirectory=/home/ec2-user/spikelines/bot
# Node lives under nvm on the box, so use absolute paths + an explicit PATH
# (systemd services get a minimal PATH that won't include the nvm bin dir).
# Adjust the nvm version path if node is upgraded.
Environment=PATH=/home/ec2-user/.nvm/versions/node/v22.23.1/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
EnvironmentFile=/home/ec2-user/spikelines/bot/.env
ExecStart=/home/ec2-user/.nvm/versions/node/v22.23.1/bin/node /home/ec2-user/.nvm/versions/node/v22.23.1/bin/npx tsx src/index.ts
Restart=always
RestartSec=5
MemoryMax=256M
CPUQuota=40%

[Install]
WantedBy=multi-user.target
